# deploy/systemd/skills-worker.service
# /thearray/git/evo/systemd/skills-worker.service
#
# Temporal worker for the skills subsystem (R2R sync, GitHub importers,
# check-source drift detection). Replaces the fire-and-forget goroutines
# that used to live in dashboard/apiserver/skills_r2r_sync.go.
#
# Hardening directives mirror evo-worker@.service so the trust boundary
# stays uniform across the evo.target group.
[Unit]
Description=evo skills Temporal worker (R2R sync, importers, check-source)
Documentation=https://github.com/thearray/git/evo
After=network.target temporal-server.service
Wants=temporal-server.service
PartOf=evo.target

[Service]
Type=exec
User=evo
Group=evo
WorkingDirectory=/opt/evo
ExecStart=/opt/evo/bin/skills-worker
# on-failure rather than always: a misconfigured worker that log.Fatals
# at startup shouldn't spin up an unbounded restart loop.
Restart=on-failure
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Environment. SKILLS_R2R_* are read inside the worker; GITHUB_TOKEN is
# the auth for the importer activity. PLATFORM_BASE_URL is shared with
# reindex-worker. /etc/evo/skills-worker.env carries the secret-bearing
# vars; the file is mode 0600 owned by root:evo on the Pi.
EnvironmentFile=-/etc/evo/skills-worker.env
Environment=TEMPORAL_ADDRESS=localhost:7233
Environment=TEMPORAL_NAMESPACE=default
Environment=LOG_LEVEL=info

# Resource limits — modest. The worker holds a small pgxpool and at
# most MaxConcurrentActivityExecutionSize=8 in-flight HTTP calls.
MemoryMax=512M
CPUQuota=50%
TasksMax=200

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/evo /var/lib/evo
UMask=0077
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
SystemCallFilter=@system-service
SystemCallArchitectures=native
RestrictSUIDSGID=true
LimitNOFILE=65535

StandardOutput=journal
StandardError=journal
SyslogIdentifier=skills-worker

[Install]
WantedBy=evo.target
